# bin/anime-notifier.service
[Unit]
Description=Anime Notifier

[Service]
# User=kurozero
WorkingDirectory=/home/kurozero/anime-notifier
# Path to executable.
ExecStart=/home/kurozero/anime-notifier/src/notifier.js
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
