# bin/anime-notifier.service
[Unit]
Description=Anime Notifier

[Service]
# User=kurozero
WorkingDirectory=/home/kurozero/anime-notifier
Environment=NODE_ENV=production
ExecStart=/home/kurozero/anime-notifier/src/notifier.ts
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
